# etc/systemd/system/restic-backup.timer
[Unit]
Description=Backup with restic on schedule

[Timer]
OnCalendar=*-*-* 02:15:00 CET
Persistent=true

[Install]
WantedBy=timers.target
